# examples/files/tmp/service/without_log/run
#!/bin/sh
exec 2>&1
set -e
exec envdir ./env \
  sh -c '
    exec vmstat 1 2>&1
  '
